# packages/atlassian-jira-software/src/main/etc/init.d/jira
!/bin/sh -e
# JIRA startup script (https://confluence.atlassian.com/jira/starting-jira-automatically-on-linux-211060709.html)
#
# chkconfig: 2345 80 20
# description: JIRA

APP=jira
export JIRA_HOME=/home/jira/jira-data
export CATALINA_PID=/var/run/jira.pid

export JAVA_HOME="/usr/java/default"
export JAVA_OPTS="-Xms1024m -Xmx2048m -Djava.awt.headless=true -Dfile.encoding=UTF-8 -XX:+UseG1GC -server ${JAVA_OPTS}"

# Location of application's bin directory
BASE=/home/jira/jira

case "$1" in
  start)
    echo "Starting $APP"

    # touch $CATALINA_PID
    # chown jira:jira $CATALINA_PID

    $BASE/bin/startup.sh &> /dev/null
    ;;

  stop)
    echo "Stopping $APP"
    $BASE/bin/shutdown.sh &> /dev/null
    echo "$APP stopped successfully"
    ;;

  stop-wait)
    echo "Stopping $APP, waiting for actual termination"
    COUNTER=1
    while [ $COUNTER -gt 0 ]; do
        $BASE/bin/shutdown.sh &> /dev/null
        sleep 1
        COUNTER=$(ps aux | grep "[c]atalina.base=$BASE" | wc -l)
        done
    echo "$APP stopped successfully"
    ;;

  restart)
     $0 stop-wait
     $0 start
     ;;

  status)
     PROCESS=$(ps aux | grep "[c]atalina.base=$BASE")
     PID=`echo $PROCESS | awk {'print $2'}`
     [ $(echo $PROCESS | wc -l) -gt 0 ] && echo "$APP (pid  $PID) is running..."
     ;;

  *)
    echo "Usage: /etc/init.d/$APP {start|stop|stop-wait|restart|status}"
    exit 1
    ;;
esac
exit 0
